# deploy/ufw-okboy-cleanup.service
# UFW OkBoy - Stale Rule Cleanup Service
# Runs as a oneshot triggered by ufw-okboy-cleanup.timer

[Unit]
Description=UFW OkBoy - Cleanup stale firewall rules

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/ufw-okboy/server
ExecStart=/opt/ufw-okboy/venv/bin/python app.py -c config.yaml cleanup --max-age 7
