# deploy/ufw-okboy-cleanup.timer
# UFW OkBoy - Daily cleanup timer
# Removes firewall rules for users who haven't knocked in 7+ days
#
# Install:
#   cp ufw-okboy-cleanup.{service,timer} /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now ufw-okboy-cleanup.timer

[Unit]
Description=Daily cleanup of stale UFW OkBoy rules

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
